# supervisord.conf
[supervisord]
nodaemon=true

[program:backend]
command=python app.py
autostart=true
autorestart=true

[program:frontend]
command=streamlit run ui.py
autostart=true
autorestart=true
